# deploy/systemd/flusso-recovery.timer
[Unit]
Description=Periodically recover interrupted Flusso jobs and A2A conversations

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=flusso-recovery.service

[Install]
WantedBy=timers.target
