# deploy/systemd/flusso-a2a-health.timer
[Unit]
Description=Periodically verify Flusso A2A runtime readiness

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=flusso-a2a-health.service

[Install]
WantedBy=timers.target
