# deploy/systemd/flusso-a2a-health.service
[Unit]
Description=Verify and repair Flusso A2A runtime readiness
Wants=network-online.target flusso-a2a.service
After=network-online.target flusso-a2a.service

[Service]
Type=oneshot
User=flusso
Group=flusso
WorkingDirectory=/opt/flusso
Environment=HOME=/home/flusso
Environment=PATH=/home/flusso/.npm-global/bin:/home/flusso/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/etc/flusso/flusso.env
ExecStart=/usr/bin/env node /opt/flusso/scripts/ensure-a2a-runtime.mjs
TimeoutStartSec=600
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
UMask=0077
